# 9.6/crontab
# Wal-e crons for backup and clean old backups
0 0 0 * *  /usr/bin/envdir /etc/wal-e.d/env /usr/local/bin/backup-wale
0 0 1 * *  /usr/bin/envdir /etc/wal-e.d/env /usr/local/bin/clean-wale
